# orientdb/supervisord/orientdb.conf
[program:orientdb]
command=/opt/orientdb/bin/dserver.sh
autostart=true
autorestart=true
stopasgroup=true