# orientdb/supervisor.d/orientdb.conf
[program:orientdb]
command=/opt/orientdb.sh
autostart=true
autorestart=true
stopasgroup=true